# systemd/eddn.target
[Unit]
Description=EDDN Services
Requires=eddn@relay.service eddn@monitor.service eddn@gateway.service

[Install]
WantedBy=multi-user.target
